# supervisord.conf
[program:newrelic]
command=/home/dotcloud/current/newrelic-daemon.x64 -f -c /home/dotcloud/current/newrelic.cfg
